# releases/flexiqueue-hosting-folder-20260313-155030/scripts/pi/flexiqueue-queue.service
# Systemd unit for Laravel queue worker (TTS generation and other queued jobs).
# Copy to the Pi: sudo cp flexiqueue-queue.service /etc/systemd/system/
# Then: sudo systemctl daemon-reload && sudo systemctl enable --now flexiqueue-queue
# Required for token/station TTS pre-generation (GenerateTokenTtsJob, GenerateStationTtsJob).

[Unit]
Description=FlexiQueue queue worker
After=network.target php8.3-fpm.service nginx.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/flexiqueue
ExecStart=/usr/bin/php artisan queue:work --tries=3
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
